# systemd/.config/systemd/user/backup-projects.service
[Unit]
Description=Run backup-projects.sh for project backups

[Service]
Type=oneshot
# %h 会自动替换为你的家目录
ExecStart=%h/.local/bin/backup-projects.sh
